# recipes-mini/canvas-server/files/canvas-server.service
[Unit]
Description=Canvas-Server startup script

[Service]
ExecStart=/usr/bin/server

[Install]
WantedBy=multi-user.target

